# systemd/backup-job@.service
[Unit]
Description=Backup job %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/backup-system/bin/backup-job.sh %i
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
KillMode=control-group
TimeoutStopSec=2m

[Install]
WantedBy=multi-user.target
